# contrib/solaredge-exporter.service
[Unit]
Description=SolarEdge Prometheus exporter: export ModBus TCP stats gathered from Solaredge inverters to Prometheus

[Service]
User=solaredge-exporter
WorkingDirectory=~
# Write the log somewhere the service user can write to (see also LOG_PATH)
Environment=LOG_PATH=/var/log/solaredge-exporter/SolarEdge-Exporter.log
ExecStart=/path/to/solaredge-exporter
PrivateTmp=yes
NoNewPrivileges=yes
Before=prometheus.service
Wants=prometheus.service

[Install]
WantedBy=multi-user.target
